# private_dot_config/autostart/custom-xkb.desktop
[Desktop Entry]
Type=Application
Name=Custom XKB Layout
Comment=Apply custom keyboard layout for Vim navigation
Exec=/home/tsuyoshi/.local/bin/apply-custom-xkb.sh
Terminal=false
Hidden=false
X-GNOME-Autostart-enabled=true
X-XFCE-Autostart-enabled=true 