# private_dot_config/systemd/user/xkb-monitor.timer
[Unit]
Description=Run XKB layout monitor every 5 minutes
Documentation=https://github.com/tsuyoshi/dotfiles

[Timer]
OnBootSec=30s
OnUnitActiveSec=5min
AccuracySec=1s

[Install]
WantedBy=timers.target 