# etc/systemd/system/amd-s2idle-optimize.service
# AMD Ryzen s2idle 電力最適化
# USB autosuspend / PCIe ASPM / NVMe 省電力を有効化して s2idle 中の消費電力を削減
[Unit]
Description=AMD s2idle power optimization for ThinkPad E14 Gen6
After=multi-user.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash -c '\
  echo 1 > /sys/module/snd_hda_intel/parameters/power_save 2>/dev/null || true; \
  for dev in /sys/bus/usb/devices/*/power/autosuspend_delay_ms; do echo 1000 > "$dev" 2>/dev/null || true; done; \
  for dev in /sys/bus/usb/devices/*/power/control; do echo auto > "$dev" 2>/dev/null || true; done; \
  for dev in /sys/bus/pci/devices/*/power/control; do echo auto > "$dev" 2>/dev/null || true; done; \
  echo powersupersave > /sys/module/pcie_aspm/parameters/policy 2>/dev/null || true; \
  echo 1 > /sys/devices/system/cpu/cpufreq/conservative/ignore_nice_load 2>/dev/null || true; \
'

[Install]
WantedBy=multi-user.target
